# etc/opennms.service
[Unit]
Description=OpenNMS server
After=postgresql.target

[Service]
Type=forking
User=root
PIDFile=/opt/opennms/logs/daemon/opennms.pid
ExecStart=/opt/opennms/bin/opennms -Q start
ExecStop=/opt/opennms/bin/opennms stop

[Install]
WantedBy=multi-user.target
