# etc/opennms.service
[Unit]
Description=OpenNMS server
Wants=postgresql.service
Requires=network.target network-online.target
After=postgresql.service postgresql-10.service postgresql-11.service postgresql-12.service postgresql-13.service network.target network-online.target

[Service]
User=root
Environment="OPENNMS_HOME=/usr/share/opennms"

Type=forking
PIDFile=/var/run/opennms.pid

ExecStart=/usr/share/opennms/bin/opennms -s start
ExecStop=/usr/share/opennms/bin/opennms stop

[Install]
WantedBy=multi-user.target
